# projects/backend-servidor/streaming/janus/janus.service
# ============================================
# JANUS GATEWAY - SYSTEMD SERVICE
# ============================================
# Target location on Linux: /etc/systemd/system/janus.service
# ============================================

[Unit]
Description=Janus WebRTC Gateway
Documentation=https://janus.conf.meetecho.com/docs/
After=network.target

[Service]
Type=simple
ExecStart=/opt/janus/bin/janus -o
Restart=always
RestartSec=5
User=root
WorkingDirectory=/opt/janus
StandardOutput=journal
StandardError=journal

# Limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target
